# systemd/nvml/fair_cluster_nvml_resources.slice
[Unit]
Description=Slice for nvml_monitor

[Slice]
MemoryHigh=700M
MemoryMax=1G
CPUQuota=100%
IOWeight=100
IODeviceWeight=/dev/sda2 100
StartupDeviceWeight=/dev/sda2 100
