# systemd/hc_resources.slice
[Unit]
Description=Slice for gcm services

[Slice]
MemoryHigh=700M
MemoryMax=1G
CPUQuota=200%
